# config/claude/scripts/com.user.weekly-review.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.weekly-review</string>
    <key>ProgramArguments</key>
    <array>
        <string>~/.claude/scripts/weekly-review.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/tmp/weekly-review.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/weekly-review.error.log</string>
    <key>WorkingDirectory</key>
    <string>~/Documents/MyLife</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:~/.local/bin</string>
        <key>HOME</key>
        <string>~</string>
    </dict>
</dict>
</plist>
